# crontab
* * * * * root cd /var/www/asteroid-neo && python state.py

